# supervisord.conf
[supervisord]
nodaemon=true

[program:mongo]
command=mongod

[program:app]
command=/bin/bash -c "python3 -u /app/app.py"
